# supervisord.conf
; =============================================================================
; TEMPEST STACK — supervisord configuration
; Manages T3MP3ST (Node.js)
; =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── Nginx reverse proxy ──────────────────────────────────────────────────────
; Terminates HTTPS on :8443 and proxies to T3MP3ST on :3333.
; Starts before t3mp3st (priority=5) so the HTTPS port is ready quickly.
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
directory=/
autostart=true
autorestart=true
priority=5
startsecs=2
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── T3MP3ST server ───────────────────────────────────────────────────────────
; Node.js server serving the T3MP3ST Web UI and REST API on port 3333.
[program:t3mp3st]
command=node /opt/t3mp3st/dist/server.js
directory=/opt/t3mp3st
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    NODE_ENV="production",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    OPENROUTER_API_KEY="%(ENV_OPENROUTER_API_KEY)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    LLM_PROVIDER="%(ENV_LLM_PROVIDER)s",
    LLM_MODEL="%(ENV_LLM_MODEL)s"
